# daemon/sunday-code-daemon/launchd/com.trentontaylor.sunday-code-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trentontaylor.sunday-code-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/trentontaylor/.nvm/versions/node/v22.21.1/bin/node</string>
    <string>/Users/trentontaylor/WebstormProjects/taylorurl-com/daemon/sunday-code-daemon/src/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/trentontaylor/WebstormProjects/taylorurl-com/daemon/sunday-code-daemon</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/trentontaylor/.nvm/versions/node/v22.21.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/trentontaylor/.sunday-code-daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/trentontaylor/.sunday-code-daemon.err.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
